# systemd/travelpi-dns-monitor.service
[Unit]
Description=TravelPi DNS Monitor
After=network.target dnsmasq.service

[Service]
Type=simple
ExecStart=/usr/local/bin/travelpi-dns-monitor.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
